# cookbooks/vine_web/templates/default/supervisord_gunicorn.conf.erb
[program:gunicorn]
command=<%= node['vine_web']['vine_venv_dir'] %>/bin/gunicorn flask_app:app -b localhost:8000 --workers=4
directory=<%= node['vine_web']['vine_repo_dir'] %>
stdout_logfile=<%= @log_dir %>/gunicorn.log
stderr_logfile=<%= @log_dir %>/gunicorn_error.log
user=<%= @env_data['server']['user'] %>
autostart=true
autorestart=true
redirect_stderr=True
priority = 1
